# qlog-client/etc/cron.d/pf-xivo-qlog-client
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# send new queue log data daily, at 23:32
32 23 * * * root    qlogc -q
